# scripts/chemist-studio.service
[Unit]
Description=Chemist Agent LangGraph Studio Dev Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/projects/chemist-agent
EnvironmentFile=/opt/projects/chemist-agent/.env
ExecStart=/opt/projects/chemist-agent/venv/bin/langgraph dev --host 0.0.0.0 --port 2024 --no-browser
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
